# autostart/.config/autostart/fusuma.sh.desktop
[Desktop Entry]
Exec=/home/sunit/Scripts/fusuma.sh
Icon=application-x-shellscript
Name=fusuma.sh
Type=Application
X-KDE-AutostartScript=true
